# raspberry-pi-5/crontab
# /var/spool/cron/crontabs/root (tai /etc/crontab)
#
# Tekijä: Jesse Karnavaara
# Vuosi: 2026
# Tarkoitus: Automatisoida Codebased Oy:n kehitysinfran (Raspberry Pi 5) 
#            Tailscale TLS-sertifikaattien elinkaarenhallinta ja 
#            3–2–1-varmuuskopiointikäytäntö.
#
# Huomiot konfiguraatiosta ja arkkitehtuurista:
# - Aikataulutuksen porrastus: Raskaat toimenpiteet on porrastettu (klo 02:00, 
#   03:00 ja 03:30) Raspberry Pi 5:n levy-I/O:n ja prosessorin ylikuormittumisen 
#   välttämiseksi. Paikallisella varmuuskopiolla on hyvin aikaa valmistua 
#   ennen offsite-siirron alkamista.
# - Luotettavuus (Absoluuttiset polut): Cron-ympäristön rajallisen PATH-muuttujan vuoksi 
#   komennot on määritetty täydellisillä absoluuttisilla poluilla (/usr/bin/...).
# - Ehdollinen logiikka (&&): NGINX-käänteisvälityspalvelin ladataan uudelleen (reload) 
#   vain, jos sertifikaatin uusimisprosessi päättyy onnistuneesti (exit code 0).
# - Varmuuskopioinnin lokitus: Varmuuskopioskriptien tulosteet ja mahdolliset virheet 
#   (stdout ja stderr, 2>&1) ohjataan suoraan omiin lokitiedostoihinsa (/var/log/...). 
#   Tämä estää virheiden katoamisen ja paikallisen MTA:n (Mail Transfer Agent) tukkeutumisen.
#
# Verifiointi:
#   Tarkista cron-loki suoritusten jälkeen: grep CRON /var/log/syslog
#   Tarkista sertifikaatin voimassaolo: tailscale cert --status
#   Tarkista varmuuskopioinnin onnistuminen: tail -n 20 /var/log/offsite_backup.log
#
# HUOM: Tailscale-verkko-osoite on osittain sensuroitu.

# M H DOM MON DOW COMMAND

# 1. TLS-sertifikaattien automaattinen uusinta ja NGINX reload (Maanantaisin klo 03:00)
0 3 * * 1 /usr/bin/tailscale cert --cert-file /etc/ssl/certs/tailscale.crt --key-file /etc/ssl/private/tailscale.key codebased-edge.[SENSUROITU].ts.net && /usr/bin/docker exec nginx_lb nginx -s reload

# 2. Paikallinen varmuuskopiointi (Joka yö klo 02:00)
# Pakkaa Dockerin nimetyt volyymit ja asetustiedostot paikalliselle levylle.
0 2 * * * /usr/local/bin/local_backup.sh >> /var/log/local_backup.log 2>&1

# 3. Offsite-varmuuskopiointi (Joka yö klo 03:30)
# Siirtää pakatut varmuuskopiot ulkoiseen sijaintiin (Hetzner Storage Box).
30 3 * * * /usr/local/bin/offsite_backup.sh >> /var/log/offsite_backup.log 2>&1
